# resources/linux/screen-recorder-picker.service
[Unit]
Description=Screen Recorder Portal Picker Service
Documentation=https://github.com/your-org/screen-recorder
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/screen-recorder-picker
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=screen_recorder_picker=debug

[Install]
WantedBy=graphical-session.target
